# mrq/supervisord_templates/default.conf
[program:mrq-worker]

command={{ SUPERVISORD_COMMAND }}
process_name=%(program_name)s%(process_num)s
numprocs={{ SUPERVISORD_PROCESSES }}
startsecs=1
startretries=3
stopsignal=TERM
autostart=true
autorestart=unexpected
exitcodes=0,2,3
stopasgroup=false
killasgroup=false

[supervisord]
nodaemon=true

logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=debug
